# components/g8ep/scripts/supervisord.conf
; g8e node supervisor configuration.
;
; Variables are interpolated at startup from the container environment using
; supervisor's native %(ENV_VAR)s syntax. The g8ep entrypoint is responsible
; for ensuring G8E_INTERNAL_AUTH_TOKEN is exported (loaded from the g8es-ssl
; volume) and G8E_SUPERVISOR_PORT is set (defaulted to 443) before invoking
; supervisord with this config.

[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/var/log/supervisor

[unix_http_server]
file=/tmp/supervisor.sock

[inet_http_server]
port = *:%(ENV_G8E_SUPERVISOR_PORT)s
username = g8e-internal
password = %(ENV_G8E_INTERNAL_AUTH_TOKEN)s

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:operator]
command=/app/components/g8ep/scripts/fetch-key-and-run.sh
autostart=true
autorestart=true
startsecs=10
startretries=999
stopwaitsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/1
stderr_logfile_maxbytes=0
